# Library/LaunchAgents/com.local.sketchybar-wake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.sketchybar-wake</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>sleep 2 &amp;&amp; killall sketchybar; sleep 1; sketchybar &amp;</string>
    </array>
    <key>LaunchEvents</key>
    <dict>
        <key>com.apple.notifyd.matching</key>
        <dict>
            <key>Wake</key>
            <dict>
                <key>Notification</key>
                <string>com.apple.powermanagement.systempowerstate</string>
            </dict>
        </dict>
    </dict>
</dict>
</plist>
